# scripts/backups/canary/pod-golly-wiki-canary.timer.j2
[Unit]
Description=Timer to run the pod-golly-wiki backups canary

[Timer]
OnCalendar=*-*-* 8:05:00

[Install]
WantedBy=timers.target
