# scripts/backups/canary/pod-golly-wiki-canary.timer.j2
[Unit]
Description=Timer to run the pod-golly-wiki backups canary

[Timer]
#OnCalendar=Sun *-*-* 9:03:00
OnCalendar=*-*-* 7:03:00

[Install]
WantedBy=timers.target
