# scripts/backups/pod-golly-wiki-backups-cleanolderthan.timer.j2
[Unit]
Description=Timer to clean files older than N days from the pod-golly-wiki backups dir

[Timer]
OnCalendar=Sun *-*-* 3:06:00

[Install]
WantedBy=timers.target
