# scripts/backups/pod-golly-wiki-backups-wikidb.timer.j2
[Unit]
Description=Timer to back up the pod-golly-wiki wiki database

[Timer]
#OnCalendar=Sun *-*-* 2:02:00
OnCalendar=*-*-* 2:02:00

[Install]
WantedBy=timers.target
