# scripts/backups/pod-golly-wiki-backups-wikifiles.timer.j2
[Unit]
Description=Timer to back up pod-golly-wiki files

[Timer]
OnCalendar=Sun *-*-* 2:08:00

[Install]
WantedBy=timers.target
