# scripts/backups/pod-golly-wiki-backups-aws.timer.j2
[Unit]
Description=Timer to copy the lastest pod-golly-wiki backup to an S3 bucket

[Timer]
OnCalendar=Sun *-*-* 2:18:00

[Install]
WantedBy=timers.target
