# scripts/pod-golly-wiki.service.j2
[Unit]
Description=Golly wiki docker pod
Requires=docker.service
After=docker.service

[Service]
Restart=always
StandardError=null
StandardOutput=null
ExecStartPre=/usr/bin/test -f {{ golly_wiki_pod_install_dir }}/docker-compose.yml
# `docker network prune -f` used to be here (added 2024-05-05 after a stale
# network interface caused an outage), but pruning is unsafe now that the
# `pod-golly-wiki_frontend` network is shared with the golly ai-slackbots stack:
# a prune while no pod-golly-wiki containers are attached could remove a network
# the slackbots are actively using. Instead, ensure the shared network exists
# (idempotent inspect-or-create) before compose brings the stack up. The compose
# file declares `frontend` as `external: true, name: pod-golly-wiki_frontend`.
ExecStartPre=/bin/bash -c '/usr/bin/docker network inspect pod-golly-wiki_frontend >/dev/null 2>&1 || /usr/bin/docker network create --driver bridge pod-golly-wiki_frontend'
ExecStart=/usr/local/bin/docker-compose -f {{ golly_wiki_pod_install_dir }}/docker-compose.yml up
ExecStop=/usr/local/bin/docker-compose  -f {{ golly_wiki_pod_install_dir }}/docker-compose.yml stop

[Install]
WantedBy=default.target
